# packaging/linux/polaris.service.in
[Unit]
Description=@PROJECT_DESCRIPTION@
After=graphical-session.target
PartOf=graphical-session.target
StartLimitIntervalSec=500
StartLimitBurst=5

[Service]
# Avoid starting Polaris before the desktop is fully initialized.
ExecStartPre=/bin/sleep 5
PassEnvironment=WAYLAND_DISPLAY DISPLAY XDG_CURRENT_DESKTOP XDG_SESSION_TYPE XDG_SESSION_ID XAUTHORITY XDG_RUNTIME_DIR DBUS_SESSION_BUS_ADDRESS
@POLARIS_SERVICE_START_COMMAND@
@POLARIS_SERVICE_STOP_COMMAND@
Restart=on-failure
RestartSec=5s
# Give capture, encode, and audio workers room to request realtime or
# negative-nice priority when the user manager allows it.
LimitRTPRIO=95
LimitNICE=-10

[Install]
WantedBy=xdg-desktop-autostart.target
